# deploy/systemd/log-lottery-ws.service
[Unit]
Description=log-lottery websocket server (rust)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/log-lottery/ws_server
ExecStart=/opt/log-lottery/ws_server/target/release/ws_server
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
